# aurvex-dashboard.service
[Unit]
Description=AURVEX.Ai Dashboard (Flask + SocketIO)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/trade_engine
EnvironmentFile=/root/trade_engine/.env
ExecStart=/root/trade_engine/venv/bin/python3 /root/trade_engine/app.py
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/root/trade_engine/logs/dashboard.log
StandardError=append:/root/trade_engine/logs/dashboard.log

[Install]
WantedBy=multi-user.target
