# aurvex-dashboard.service
[Unit]
Description=AURVEX.Ai Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/trade_engine
ExecStart=/root/trade_engine/venv/bin/python3 /root/trade_engine/app.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
