# aurvex-bot.service
[Unit]
Description=AURVEX.Ai Scalp Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/trade_engine
ExecStart=/root/trade_engine/venv/bin/python3 /root/trade_engine/scalp_bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
