# aurvex-dashboard.service
[Unit]
Description=AURVEX.Ai Dashboard (Flask + SocketIO)
After=network.target aurvex-bot.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/trade-engine
EnvironmentFile=-/root/trade-engine/.env
ExecStart=/root/trade-engine/venv/bin/python3 /root/trade-engine/app.py
Restart=always
RestartSec=10
StandardOutput=append:/root/trade-engine/logs/dashboard.log
StandardError=append:/root/trade-engine/logs/dashboard.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
