# aurvex-watchdog.service
[Unit]
Description=AURVEX.Ai Watchdog
After=network.target aurvex-bot.service aurvex-dashboard.service
Wants=aurvex-bot.service aurvex-dashboard.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/trade-engine
ExecStart=/root/trade-engine/venv/bin/python3 /root/trade-engine/watchdog.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
