# systemd/ax-dashboard.service
[Unit]
Description=AX Trade Engine - Dashboard API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/trade_engine
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/trade_engine/.venv/bin/python app.py
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=15
StandardOutput=append:/var/log/ax_dashboard.log
StandardError=append:/var/log/ax_dashboard.log

[Install]
WantedBy=multi-user.target
